# hyprland/.config/systemd/user/hypr-sunsetr.service
[Unit]
Description=Sunsetr - Automatic color temperature adjustment for Hyprland
PartOf=graphical-session.target
Requires=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/sunsetr
Restart=on-failure
RestartSec=30

[Install]
WantedBy=graphical-session.target
