# waybar/.config/systemd/user/setup-temps.service
[Unit]
Description=Symlink temperatures to /tmp

[Service]
Type=oneshot
ExecStart=%h/.local/bin/setup-temps.sh

[Install]
WantedBy=default.target
