# liquidctl/.config/systemd/user/liquidctl.service
[Unit]
Description=Liquidctl

[Service]
Type=oneshot
Restart=on-failure
RestartSec=1s
ExecStart=liquidctl initialize all
ExecStart=liquidctl set pump speed 90
ExecStart=liquidctl set logo color blackout

[Install]
WantedBy=default.target
